# build/dailywire-podcast-proxy.service
[Unit]
Description=DailyWire Podcast Proxy

[Service]
Type=notify
WorkingDirectory=/usr/lib/dailywire-podcast-proxy
Environment="ASPNETCORE_URLS=http://*:9473"
ExecStart=/usr/lib/dailywire-podcast-proxy/DailyWirePodcastProxy

[Install]
WantedBy=multi-user.target